# debian/crowdsec-hubupdate.service
[Unit]
Description=CrowdSec Hub update

[Service]
Type=oneshot
Environment=LC_ALL=C LANG=C PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/libexec/crowdsec/hubupdate.sh

