# config/crowdsec.cron.daily
#!/bin/sh

test -x /usr/bin/cscli || exit 0

/usr/bin/cscli --error hub update

upgraded=$(/usr/bin/cscli --error hub upgrade)
if [ -n "$upgraded" ]; then
    systemctl reload crowdsec
fi

exit 0
