# config/crowdsec.cron.daily
#!/bin/sh
test -x /usr/bin/cscli || exit 0
/usr/bin/cscli --error hub update && /usr/bin/cscli --error hub upgrade
systemctl reload crowdsec
exit 0
