# config/crowdsec.cron.daily
#!/bin/sh
/usr/bin/cscli --error hub update && /usr/bin/cscli --error hub upgrade
if [ $(/usr/bin/systemctl is-enabled crowdsec.service) = "enabled" ]; then
    /usr/bin/systemctl reload crowdsec.service
fi
exit 0