# config/crowdsec.service
[Unit]
Description=CrowdSec agent
After=syslog.target network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
PIDFile=/var/run/crowdsec.pid
ExecStart=/usr/local/bin/crowdsec -c /etc/crowdsec/config/config.yaml
ExecStartPost=/bin/sleep 0.1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
